# services/evdev_transformer@.service
# vi: ft=systemd
[Unit]
Description=evdev transformer

[Service]
ExecStart=python -m evdev_transformer %I

[Install]
WantedBy=default.target
